# scripts/agent-bus.service
[Unit]
Description=Agent Bus — Cross-machine event relay for AI agents
After=network.target

[Service]
Type=simple
User=agent-bus
Group=agent-bus
WorkingDirectory=/opt/agent-bus/app
EnvironmentFile=/etc/agent-bus/.env
ExecStart=/opt/agent-bus/app/.venv/bin/python -m uvicorn server.main:app --host ${AGENT_BUS_HOST} --port ${AGENT_BUS_PORT}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-bus

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/agent-bus/data
ReadOnlyPaths=/opt/agent-bus/app /etc/agent-bus

[Install]
WantedBy=multi-user.target
